# agent-voice-chat/vm/swarm-server.service
[Unit]
Description=AI Agents X Space (three.ws) Realtime server
After=network.target

[Service]
Type=simple
User=agent
WorkingDirectory=/home/agent/ai-agents-x-space
ExecStart=/usr/bin/node /home/agent/ai-agents-x-space/index.js
Restart=on-failure
StandardOutput=append:/home/agent/server.log
StandardError=append:/home/agent/server.log

[Install]
WantedBy=multi-user.target
